# etc/supervisord.d/redis.ini
[program:redis]
command=/usr/bin/redis-server --port 6379
autorestart=true
