# etc/supervisord.d/soteria-polling-loop.ini
[program:soteria-polling-loop]
command=flask soteria run-polling-loop
directory=/srv
autorestart=true
