# etc/supervisord.d/soteria-polling-loop.ini
[program:soteria-polling-loop]
command=flask soteria run-polling-loop
autorestart=true
user=apache
directory=/srv
